# rpi4-k8s-master/etc/systemd/system/srv-nfs-lib.automount
[Unit]
Description=Automatically mount NFS lib folder

[Mount]
Where=/srv/nfs/lib

[Install]
WantedBy=local-fs.target
