# rpi4-k8s-master/etc/systemd/system/srv-nfs-lib.mount
[Unit]
Description=Mount NFS lib folder

[Mount]
What=/var/lib
Where=/srv/nfs/lib
Options=bind
TimeoutSec=5m

[Install]
WantedBy=local-fs.target

