# rpi4-k8s-master/etc/systemd/system/srv-nfs-music.automount
[Unit]
Description=Automatically mount NFS Music folder

[Mount]
Where=/srv/nfs/music

[Install]
WantedBy=local-fs.target
