# configs/cni-dhcp-daemon.service
[Unit]
Description=CNI DHCP Daemon

[Service]
ExecStart=/opt/cni/bin/dhcp daemon

[Install]
WantedBy=multi-user.target